# redis_m1sN_make_sentinel_ka_ha/roles/redis/templates/default/sentinel.service.j2
[Unit]
Description=Redis Sentinel
After=network.target

[Service]
User=redis
Group=redis
ExecStart={{ redis_install_dir }}/bin/redis-server /etc/redis/sentinel_{{ sentinel_port }}.conf --sentinel
Restart=always

[Install]
WantedBy=multi-user.target
